# backend/portfolio-backend.service
[Unit]
Description=Portfolio Backend API - Gunicorn
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/www/portfolio/backend
Environment=PATH="/var/www/portfolio/backend/.venv/bin"
ExecStart=/var/www/portfolio/backend/.venv/bin/gunicorn main:app -w 4 -k uvicorn.workers.UvicornWorker --bind unix:/run/portfolio.sock --timeout 120 --access-logfile - --error-logfile -
Restart=always
RestartSec=5
PrivateTmp=true

[Install]
WantedBy=multi-user.target
